# docker/config/supervisor/conf.d/redis.conf
[program:redis]
chown=redis-og:redis-og
user=redis-og
command=/usr/local/bin/redis-server
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0